# supervisor/conf.d/engine-001.conf
[program:engine-001]
directory=/opt/ffplayout-engine
command=./venv/bin/python3 ffplayout.py -c /etc/ffplayout/ffplayout-001.yml
killasgroup=true
stopasgroup=true
autorestart=true
autostart=true
startsecs=2
startretries=10
