# supervisor/conf.d/engine-001.conf
[program:engine-001]
directory=/opt/ffplayout_engine
command=/opt/ffplayout_engine/venv/bin/python ffplayout.py -c /etc/ffplayout/ffplayout-001.yml
redirect_stderr=true
stdout_logfile=/var/log/ffplayout/engine-001.log
killasgroup=true
stopasgroup=true
autorestart=true
autostart=true
startsecs=2
startretries=10
