# supervisor/conf.d/engine-001.conf
[program:ffplayout-engine-001]
directory=/opt/ffplayout-engine
command=./venv/bin/python3 ffplayout.py # -c ffplayout.yml
killasgroup=true
stopasgroup=true
autorestart=true
autostart=true
startsecs=2
startretries=10
